# roles/firewall/templates/ansible-pull/ansible-pull-reconfigure.path
[Path]
PathModified={{ ansible_inventory }}

#[Install]
#WantedBy=multi-user.target
