# roles/firewall/templates/ansible-pull/ansible-pull-reconfigure.path
[Path]
PathModified={{ firewall_ansible_inventory }}

[Install]
WantedBy=multi-user.target
