# roles/firewall/files/ansible-pull/ansible-pull.timer
[Unit]
Description=Run ansible-pull on a timer.

[Timer]
OnUnitActiveSec=60min

[Install]
WantedBy=timers.target
