# roles/base/files/ansible-pull-reconfigure.path
[Path]
PathModified=/boot/firmware/inventory.yml

[Install]
WantedBy=multi-user.target
