# roles/base/files/ansible-pull.timer
[Unit]
Description=Run ansible-pull on a timer.

[Timer]
OnBootSec=1min
OnUnitActiveSec=5min

[Install]
WantedBy=timers.target
